# scripts/crontab
*/7 * * * * cd $HOME/mirage-tutorial-deployment && git pull
